# docker/conf/nodespeed-ide.conf
[program:nodespeed-ide]
command = /start.sh -d -s /projects/.brackets-server -j /projects
autostart = true
autorestart = true
environment = NODE_ENV="dev"
redirect_stderr=true

